# salt/ospf/zebra.service
[Unit]
After = network.target

[Service]
ExecStart = /usr/lib/quagga/zebra

[Install]
WantedBy = default.target
